# deploy/systemd/harbor-relay-agent.service
[Unit]
Description=Harbor Relay Agent
Documentation=https://github.com/yuanyp8/harbor-relay
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/harbor-relay-agent
ExecStart=/usr/local/bin/harbor-relay-agent --config /etc/harbor-relay/agent.yaml
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
